# common/debian/contrail-database-common/debian/contrail-database-common.postinst
#!/bin/sh

echo "Running Postinst for contrail-database-common install.."
    sudo update-rc.d cassandra disable
    sudo update-rc.d contrail-database defaults
echo "Postinst for contrail-openstack-database-common done"
