# common/debian/contrail-database-common/debian/contrail-database-common.postinst
#!/bin/sh

if [ "$1" = "configure" ]; then
    echo "Running Postinst for contrail-database-common install.."
    sudo update-rc.d cassandra disable
    sudo update-rc.d contrail-database defaults
fi
for svc in cassandra ; do
    if [ -f /etc/init/$svc.override ]; then
            rm /etc/init/$svc.override
    fi
done
echo "Postinst for contrail-database-common done"
