# storage/control_files/storage-stats.upstart
description "Storage statistics daemon"
author "Shreya Kunjir <shreyak@juniper.net>"

start on runlevel [2345] or storage-stats
stop on runlevel [!2345]

respawn
respawn limit 5 30

exec /usr/bin/contrail-storage-stats --conf_file /etc/contrail/contrail-storage-nodemgr.conf
