# common/control_files/barbican-api.initd.supervisord
#!/bin/sh

# chkconfig: 2345 99 01
# description: Juniper Network Virtualization API

if [ -e /tmp/supervisord_openstack.sock ]; then
    supervisorctl -s unix:///tmp/supervisord_openstack.sock ${1} `basename ${0}`
else
    supervisorctl -s unix:///var/run/supervisord_openstack.sock ${1} `basename ${0}`
