# common/control_files/barbican-api.ini

[program:barbican-api]
command=/bin/bash -c "exec start-stop-daemon --start --chuid barbican --chdir /var/lib/barbican --name barbican-api --exec /usr/bin/barbican-api"
redirect_stderr=true
stdout_logfile=/var/log/contrail/barbican-api-supervisor-stdout.log
stderr_logfile=/var/log/contrail/barbican-api-supervisor-stderr.log
priority=440
autostart=true
killasgroup=true
stopasgroup=true
stopsignal=KILL
exitcodes=0
