# common/control_files/barbican-worker.ini

[program:barbican-worker]
command=/bin/bash -c "exec start-stop-daemon --start --chuid barbican --chdir /var/lib/barbican --name barbican-worker --exec /usr/bin/barbican-worker -- --config-file=/etc/barbican/barbican.conf --logto /var/log/barbican/barbican-worker.log"
redirect_stderr=true
stdout_logfile=/var/log/contrail/barbican-worker-supervisor-stdout.log
stderr_logfile=/var/log/contrail/barbican-worker-supervisor-stderr.log
priority=440
autostart=true
killasgroup=true
stopasgroup=true
stopsignal=KILL
exitcodes=0
