# common/control_files/ceilometer-alarm-evaluator.ini

[program:ceilometer-alarm-evaluator]
command=/bin/bash -c "exec start-stop-daemon --start --chuid ceilometer --exec /usr/bin/ceilometer-alarm-evaluator -- --log-dir=/var/log/ceilometer"
redirect_stderr=true
stdout_logfile=/var/log/contrail/ceilometer-alarm-evaluator-supervisor-stdout.log
stderr_logfile=/var/log/contrail/ceilometer-alarm-evaluator-supervisor-stderr.log
priority=440
autostart=true
killasgroup=true
stopasgroup=true
stopsignal=KILL
exitcodes=0
