# common/control_files/ceilometer-agent-central.ini

[program:ceilometer-agent-central]
command=/bin/bash -c "exec start-stop-daemon --start --chuid ceilometer --exec /usr/bin/ceilometer-agent-central"
redirect_stderr=true
stdout_logfile=/var/log/contrail/ceilometer-agent-central-supervisor-stdout.log
stderr_logfile=/var/log/contrail/ceilometer-agent-central-supervisor-stderr.log
priority=440
autostart=true
killasgroup=true
stopasgroup=true
stopsignal=KILL
exitcodes=0
