# common/control_files/ceilometer-api.ini

[program:ceilometer-api]
command=/bin/bash -c "exec start-stop-daemon --start --chuid ceilometer --exec /usr/bin/ceilometer-api -- --log-dir=/var/log/ceilometer"
redirect_stderr=true
stdout_logfile=/var/log/contrail/ceilometer-api-supervisor-stdout.log
stderr_logfile=/var/log/contrail/ceilometer-api-supervisor-stderr.log
priority=440
autostart=true
killasgroup=true
stopasgroup=true
stopsignal=KILL
exitcodes=0
