# common/control_files/ceilometer-agent-notification.ini

[program:ceilometer-agent-notification]
command=/bin/bash -c "exec start-stop-daemon --start --chuid ceilometer --exec /usr/bin/ceilometer-agent-notification"
redirect_stderr=true
stdout_logfile=/var/log/contrail/ceilometer-agent-notification-supervisor-stdout.log
stderr_logfile=/var/log/contrail/ceilometer-agent-notification-supervisor-stderr.log
priority=440
autostart=true
killasgroup=true
stopasgroup=true
stopsignal=KILL
exitcodes=0
