# common/control_files/ceilometer-alarm-notifier.ini

[program:ceilometer-alarm-notifier]
command=/bin/bash -c "exec start-stop-daemon --start --chuid ceilometer --exec /usr/bin/ceilometer-alarm-notifier -- --log-dir=/var/log/ceilometer"
redirect_stderr=true
stdout_logfile=/var/log/contrail/ceilometer-alarm-notifier-supervisor-stdout.log
stderr_logfile=/var/log/contrail/ceilometer-alarm-notifier-supervisor-stderr.log
priority=440
autostart=true
killasgroup=true
stopasgroup=true
stopsignal=KILL
exitcodes=0
