# common/control_files/cinder-scheduler.ini

[program:cinder-scheduler]
command=/bin/bash -c "exec start-stop-daemon --start --chuid cinder --exec /usr/bin/cinder-scheduler -- --config-file=/etc/cinder/cinder.conf --log-file=/var/log/cinder/cinder-scheduler.log"
redirect_stderr=true
priority=440
autostart=true
killasgroup=true
stopasgroup=true
stopsignal=KILL
exitcodes=0
