# common/control_files/contrail-svc-monitor.initd.supervisord
#!/usr/bin/env bash

# chkconfig: 2345 99 01
# description: Juniper Network Virtualization Svc-monitor

supervisorctl -s http://localhost:9004 ${1} `basename ${0}`
