# common/control_files/contrail-dns.initd.supervisord
#!/bin/sh

# chkconfig: 2345 99 01
# description: Juniper Network Virtualization Named

supervisorctl -s http://localhost:9006 ${1} `basename ${0}`
