# common/control_files/contrail-dns.ini
[program:contrail-dns]
command= bash -c "/usr/bin/dnsd --map-server-url https://${IFMAP_SERVER}:${IFMAP_PORT} --map-user ${IFMAP_USER} --map-password ${IFMAP_PASWD} --host-ip ${HOSTIP} --discovery-server ${DISCOVERY} ${CERT_OPTS} ${LOG_LOCAL}"
priority=420
autostart=true
killasgroup=true
stopsignal=KILL
stdout_capture_maxbytes=1MB 
environment_file= /etc/contrail/dns_param
redirect_stderr=true
stdout_logfile=/var/log/contrail/dns.log
stderr_logfile=/dev/null
startsecs=5
