# common/control_files/contrail-hamon.upstart
description "contrail cmon checker"

start on (local-filesystems and net-device-up IFACE!=lo)
stop on runlevel [016]
limit core unlimited unlimited

limit nofile 8192 8192

pre-start script
  logger "Starting contrail-hamon"
end script

exec /opt/contrail/bin/contrail-ha-check.sh

respawn

post-start script
  logger "contrail-hamon started"
end script
