# common/control_files/contrail-schema-centos.ini
[program:contrail-schema]
command=/bin/bash -c "source /opt/contrail/api-venv/bin/activate && exec python /opt/contrail/api-venv/lib/python2.7/site-packages/schema_transformer/to_bgp.py --conf_file /etc/contrail/schema_transformer.conf"
priority=450
autostart=true
autorestart=true
killasgroup=true
stopsignal=KILL
redirect_stderr=true
stdout_logfile=/var/log/contrail/contrail-schema.log
stderr_logfile=/dev/null
exitcodes=0                   ; 'expected' exit codes for process (default 0,2)
