# common/control_files/contrail-opserver.ini
[program:contrail-opserver]
environment_file=/etc/contrail/opserver_param
command=bash -c "python /opt/contrail/analytics-venv/lib/python2.7/site-packages/opserver/opserver.py --redis_server_port ${REDIS_SERVER_PORT} --redis_query_port ${REDIS_QUERY_PORT} --host_ip ${HOST_IP} --http_server_port ${HTTP_SERVER_PORT} --rest_api_port ${REST_API_PORT} --disc_server_ip ${DISCOVERY}"
priority=440
autostart=true
killasgroup=true
stopsignal=KILL
stdout_capture_maxbytes=1MB 
redirect_stderr=true
stdout_logfile=/var/log/contrail/opserver.log
stderr_logfile=/dev/null
startsecs=5
