# common/control_files/contrail-collector.ini
[program:contrail-collector]
command= bash -c "/usr/bin/vizd --cassandra-server-list ${CASSANDRA_SERVER_LIST} --discovery-server ${DISCOVERY} --host-ip ${HOST_IP} --listen-port ${LISTEN_PORT} --http-server-port ${HTTP_SERVER_PORT} --analytics-data-ttl ${ANALYTICS_DATA_TTL} ${LOG_LOCAL} ${ANALYTICS_SYSLOG_PORT}"
priority=420
autostart=true
killasgroup=true
stopsignal=KILL
stdout_capture_maxbytes=1MB 
environment=pre_exec_cmd=/usr/bin/contrail_collector_pre
environment_file= /etc/contrail/vizd_param
redirect_stderr=true
stdout_logfile=/var/log/contrail/collector.log
stderr_logfile=/dev/null
startsecs=5
