# common/control_files/contrail-hamon.upstart
description "contrail cmon checker"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

limit nofile 8192 8192

exec /opt/contrail/bin/contrail-ha-check.sh
