# common/control_files/contrail-zookeeper.ini
[program:contrail-zookeeper]
command=zkServer.sh start-foreground
priority=405
autostart=true
autorestart=true
killasgroup=true
stopsignal=KILL
user=zookeeper
serverurl=http://localhost:9004
redirect_stderr=true
stdout_logfile=/var/log/contrail/zookeeper.log
stderr_logfile=/dev/null
startsecs=5
