# common/control_files/contrail-named.initd.supervisord
#!/bin/sh

# chkconfig: 2345 99 01
# description: Juniper Network Virtualization Named
ulimit -c unlimited
supervisorctl -s http://localhost:9003 ${1} `basename ${0}`
