# common/control_files/contrail-qe.initd.supervisord
#!/usr/bin/env bash

# chkconfig: 2345 99 01
# description: Juniper Network Virtualization QE
ulimit -c unlimited
supervisorctl -s http://localhost:9002 ${1} `basename ${0}`
