# common/control_files/contrail-vrouter.initd.supervisord
#!/bin/sh

# chkconfig: 2345 11 01
# description: Juniper Network Virtualization service offering

supervisorctl -s http://localhost:9005 ${1} `basename ${0}`
