# common/control_files/contrail-qe.ini
[program:contrail-qe]
command=/usr/bin/qed
priority=430
autostart=true
killasgroup=true
stopsignal=KILL
environment=pre_exec_cmd=/usr/bin/contrail_qe_pre
environment_file=/etc/contrail/qe_param
stdout_capture_maxbytes=1MB 
redirect_stderr=true
stdout_logfile=/var/log/contrail/qe.log
stderr_logfile=/dev/null
startsecs=5
