# common/control_files/glance-api.ini

[program:glance-api]
command=/bin/bash -c "exec start-stop-daemon --start --chuid glance --chdir /var/lib/glance --name glance-api --exec /usr/bin/glance-api"
redirect_stderr=true
stdout_logfile=/var/log/contrail/glance-api-supervisor-stdout.log
stderr_logfile=/var/log/contrail/glance-api-supervisor-stderr.log
priority=440
autostart=true
killasgroup=true
stopasgroup=true
stopsignal=KILL
exitcodes=0
