# common/control_files/glance-registry.ini

[program:glance-registry]
command=/bin/bash -c "exec start-stop-daemon --start --chuid glance --chdir /var/lib/glance --name glance-registry --exec /usr/bin/glance-registry"
redirect_stderr=true
stdout_logfile=/var/log/contrail/glance-registry-supervisor-stdout.log
stderr_logfile=/var/log/contrail/glance-registry-supervisor-stderr.log
priority=440
autostart=true
killasgroup=true
stopasgroup=true
stopsignal=KILL
exitcodes=0
