# common/control_files/heat-engine.ini

[program:heat-engine]
command=/bin/bash -c "exec start-stop-daemon --start --chuid heat --exec /usr/bin/heat-engine"
user=heat
redirect_stderr=true
stdout_logfile=/var/log/contrail/heat-engine-supervisor-stdout.log
stderr_logfile=/var/log/contrail/heat-engine-supervisor-stderr.log
priority=440
autostart=true
killasgroup=true
stopasgroup=true
stopsignal=KILL
exitcodes=0
