# common/control_files/heat-api-cfn.ini

[program:heat-api-cfn]
command=/bin/bash -c "exec start-stop-daemon --start --chuid heat --exec /usr/bin/heat-api-cfn"
user=heat
redirect_stderr=true
stdout_logfile=/var/log/contrail/heat-api-cfn-supervisor-stdout.log
stderr_logfile=/var/log/contrail/heat-api-cfn-supervisor-stderr.log
priority=440
autostart=true
killasgroup=true
stopasgroup=true
stopsignal=KILL
exitcodes=0
