# common/control_files/keystone.ini

[program:keystone]
command=/bin/bash -c "exec start-stop-daemon --start --chuid keystone --chdir /var/lib/keystone --name keystone --exec /usr/bin/keystone-all"
redirect_stderr=true
stdout_logfile=/var/log/contrail/keystone-supervisor-stdout.log
stderr_logfile=/var/log/contrail/keystone-supervisor-stderr.log
priority=440
autostart=true
killasgroup=true
stopasgroup=true
stopsignal=KILL
exitcodes=0
