# common/control_files/nova-scheduler.ini

[program:nova-scheduler]
command=/bin/bash -c "exec start-stop-daemon --start --chuid nova --exec /usr/bin/nova-scheduler -- --config-file=/etc/nova/nova.conf"
redirect_stderr=true
stdout_logfile=/var/log/contrail/nova-scheduler-supervisor-stdout.log
stderr_logfile=/var/log/contrail/nova-scheduler-supervisor-stderr.log
priority=440
autostart=true
startretries=10
killasgroup=true
stopasgroup=true
stopsignal=KILL
exitcodes=0
