# common/control_files/nova-console.ini

[program:nova-console]
command=/bin/bash -c "exec start-stop-daemon --start --chuid nova --exec /usr/bin/nova-console -- --config-file=/etc/nova/nova.conf"
redirect_stderr=true
stdout_logfile=/var/log/contrail/nova-console-supervisor-stdout.log
stderr_logfile=/var/log/contrail/nova-console-supervisor-stderr.log
priority=440
autostart=true
killasgroup=true
stopasgroup=true
stopsignal=KILL
exitcodes=0
