# common/control_files/nova-conductor.ini

[program:nova-conductor]
command=/bin/bash -c "exec start-stop-daemon --start --chuid nova --exec /usr/bin/nova-conductor -- --config-file=/etc/nova/nova.conf"
redirect_stderr=true
stdout_logfile=/var/log/contrail/nova-conductor-supervisor-stdout.log
stderr_logfile=/var/log/contrail/nova-conductor-supervisor-stderr.log
priority=440
autostart=true
killasgroup=true
stopasgroup=true
stopsignal=KILL
exitcodes=0
