# common/control_files/nova-consoleauth.initd.supervisord
#!/bin/sh

# chkconfig: 2345 99 01
# description: Juniper Network Virtualization API

supervisorctl -s unix:///tmp/supervisord_openstack.sock ${1} `basename ${0}`
