# common/control_files/rabbitmq-server.ini
[program:rabbitmq-server]
command=/usr/sbin/rabbitmq-server
redirect_stderr=true
priority=410
autostart=true
autorestart=true
startsecs=5
startretries=15
killasgroup=true
stopasgroup=true
stopsignal=KILL
exitcodes=0
