# common/control_files/rabbitmq-server.initd.supervisord
#!/bin/sh

# chkconfig: 2345 99 01
# description: Juniper Network Virtualization API

supervisorctl -s unix:///tmp/supervisord_support_service.sock ${1} `basename ${0}`
