# common/control_files/redis-uve.ini
[program:redis-uve]
command=redis-server /etc/contrail/redis-uve.conf
autostart=true
user=redis
priority=403
killasgroup=true
stopsignal=KILL
environment=pre_exec_cmd="/bin/rm -f /var/lib/redis/dump-uve.rdb"
stdout_capture_maxbytes=1MB
redirect_stderr=true
stdout_logfile=/var/log/contrail/redis-uve.log
stderr_logfile=/dev/null
startsecs=5
exitcodes=0                   ; 'expected' exit codes for process (default 0,2)
