# common/control_files/redis-query.ini
[program:redis-query]
command=redis-server /etc/contrail/redis-query.conf
autostart=true
user=redis
priority=402
killasgroup=true
stopsignal=KILL
stdout_capture_maxbytes=1MB 
redirect_stderr=true
stdout_logfile=/var/log/contrail/redis-query.log
stderr_logfile=/dev/null
startsecs=5
