# common/control_files/redis-sentinel.ini
[program:redis-sentinel]
command=redis-server /etc/contrail/sentinel.conf --sentinel
autostart=true
user=redis
priority=401
killasgroup=true
stopsignal=KILL
stdout_capture_maxbytes=1MB 
redirect_stderr=true
stdout_logfile=/var/log/contrail/redis-sentinel.log
stderr_logfile=/dev/null
startsecs=5
