# common/control_files/redis-sentinel.initd.supervisord
#!/bin/sh

# chkconfig: 2345 99 01
# description: Juniper Redis Sentinel

supervisorctl -s http://localhost:9002 ${1} `basename ${0}`
