# common/control_files/redis-query.initd.supervisord
#!/bin/sh

# chkconfig: 2345 99 01
# description: Juniper Redis Query

supervisorctl -s http://localhost:9002 ${1} `basename ${0}`
