# common/control_files/redis-sentinel.ini
[program:redis-sentinel]
command=/usr/sbin/redis-server /etc/contrail/sentinel.conf --sentinel
autostart=false
user=redis
priority=401
killasgroup=true
stopsignal=KILL
stdout_capture_maxbytes=1MB 
redirect_stderr=true
startsecs=5
