# common/control_files/redis-config.ini
[program:redis-config]
command=redis-server /etc/contrail/redis_config.conf
autostart=true
user=redis
priority=403
killasgroup=true
stopsignal=KILL
stdout_capture_maxbytes=1MB 
