# common/debian/contrail-ostack-dashboard/debian/grizzly/openstack-dashboard.postrm
#!/bin/sh

set -e

if [ -x "`which invoke-rc.d 2>/dev/null`" -a -x "/etc/init.d/apache2" ] ; then
	invoke-rc.d --quiet apache2 reload
fi

#DEBHELPER#
