# common/debian/contrail-openstack-database/debian/contrail-openstack-database.postinst
#!/bin/sh

if [ "$1" = "configure" ]; then
    echo "Running Postinst for contrail-openstack-database install.."
    if [ -z "$2" ]; then
        sudo service cassandra stop
        sudo rm -rf /var/lib/cassandra
    fi
    sudo update-rc.d cassandra disable
    sudo update-rc.d contrail-database defaults
fi
echo "Postinst for contrail-openstack-database done"
